# supervisor.ini
[program:prince]
directory = /opt/prince/server
command = pidproxy wsgi.pid bin/gunicorn -p wsgi.pid --workers 8 -t 120 -b :9999
user = prince
autostart = true
autorestart = true
redirect_stderr = true
